# deploy/gdex-observation-download.service
[Unit]
Description=Download 2026 UCAR GDEX observation and PREPBUFR data
Documentation=https://osdf-director.osg-htc.org/ncar/gdex/d735000/
Documentation=https://osdf-director.osg-htc.org/ncar/gdex/d337000/tarfiles/2026/
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=__APP_DIR__
EnvironmentFile=-%h/.config/gdex-observation-download.env
ExecStart=__APP_DIR__/scripts/run_daily.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
